# crontab
0 2 * * * cd /var/www && php artisan refresh:info >> /var/log/cron.log 2>&1

